# firstboot_release_gate.service
[Unit]
Description=Compose passive firstboot release-gate evidence
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/firstboot_release_gate.py --firstboot-manifest /var/log/host_vm_policy_firstboot_manifest.json --model-card /var/log/nn_ids_model_card.json --output /var/log/firstboot_release_gate.json --markdown /var/log/firstboot_release_gate.md --summary /var/log/firstboot_release_gate.summary.env --max-artifact-age-minutes 240
ExecStartPost=/bin/sh -c '/usr/local/bin/firstboot_release_gate_status.py --summary /var/log/firstboot_release_gate.summary.env --format json > /var/log/firstboot_release_gate.status.json || true'
ExecStartPost=/usr/local/bin/firstboot_release_gate_bundle_manifest.py --gate-json /var/log/firstboot_release_gate.json --gate-markdown /var/log/firstboot_release_gate.md --summary /var/log/firstboot_release_gate.summary.env --status-json /var/log/firstboot_release_gate.status.json --output /var/log/firstboot_release_gate.bundle_manifest.json
ExecStartPost=/usr/local/bin/firstboot_release_gate_bundle_manifest.py --gate-json /var/log/firstboot_release_gate.json --gate-markdown /var/log/firstboot_release_gate.md --summary /var/log/firstboot_release_gate.summary.env --status-json /var/log/firstboot_release_gate.status.json --output /var/log/firstboot_release_gate.bundle_manifest.md --format markdown
ExecStartPost=/usr/local/bin/firstboot_release_gate_operator_digest.py --status-json /var/log/firstboot_release_gate.status.json --bundle-json /var/log/firstboot_release_gate.bundle_manifest.json --output /var/log/firstboot_release_gate.operator_digest.json
ExecStartPost=/usr/local/bin/firstboot_release_gate_operator_digest.py --status-json /var/log/firstboot_release_gate.status.json --bundle-json /var/log/firstboot_release_gate.bundle_manifest.json --output /var/log/firstboot_release_gate.operator_digest.md --format markdown
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_index.py --output /var/log/firstboot_release_gate.handoff_index.json
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_index.py --output /var/log/firstboot_release_gate.handoff_index.md --format markdown
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_verify.py --index /var/log/firstboot_release_gate.handoff_index.json --artifact-root /var/log --output /var/log/firstboot_release_gate.handoff_verify.json
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_verify.py --index /var/log/firstboot_release_gate.handoff_index.json --artifact-root /var/log --output /var/log/firstboot_release_gate.handoff_verify.md --format markdown
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_freshness.py --input /var/log/firstboot_release_gate.handoff_verify.json --output /var/log/firstboot_release_gate.handoff_freshness.json --summary /var/log/firstboot_release_gate.handoff_freshness.summary.env --max-artifact-age-minutes 240
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_freshness.py --input /var/log/firstboot_release_gate.handoff_verify.json --output /var/log/firstboot_release_gate.handoff_freshness.md --format markdown --max-artifact-age-minutes 240
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_summary_smoke.py --input /var/log/firstboot_release_gate.handoff_freshness.summary.env --output /var/log/firstboot_release_gate.handoff_summary_smoke.json
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_summary_smoke.py --input /var/log/firstboot_release_gate.handoff_freshness.summary.env --output /var/log/firstboot_release_gate.handoff_summary_smoke.md --format markdown
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_status_reader.py --input /var/log/firstboot_release_gate.handoff_summary_smoke.json --format json --output /var/log/firstboot_release_gate.handoff_status_reader.json --summary /var/log/firstboot_release_gate.handoff_status_reader.summary.env
ExecStartPost=/usr/local/bin/firstboot_release_gate_handoff_status_reader.py --input /var/log/firstboot_release_gate.handoff_summary_smoke.json --format markdown --output /var/log/firstboot_release_gate.handoff_status_reader.md
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
ReadOnlyPaths=/var/log/host_vm_policy_firstboot_manifest.json /var/log/nn_ids_model_card.json
ReadWritePaths=/var/log
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
