# host_vm_policy_approval_check.timer
[Unit]
Description=Periodically validate reviewed host VM policy restore approval

[Timer]
OnBootSec=20min
OnUnitActiveSec=6h
RandomizedDelaySec=10min
Persistent=true
Unit=host_vm_policy_approval_check.service

[Install]
WantedBy=timers.target
