# host_vm_policy_attest.timer
[Unit]
Description=Refresh host/VM hardening policy attestation snapshot

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true
Unit=host_vm_policy_attest.service

[Install]
WantedBy=timers.target
