# nn_ids_report.timer
[Unit]
Description=Run IDS report hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
